# overlay/system/usr/lib/systemd/system/lxc-android-config.service.d/gsi-container-config.conf
[Service]
# Adjust LXC container configuration to skip unavailable namespaces
ExecStartPre=-/usr/libexec/lxc-android-config/gsi-android-container-config
